# profile.d/alias.sh
alias ll='ls -la'
alias gitpush='git push --force-with-lease --force-if-includes'
